# deploy/systemd/marketing-worker.service
[Unit]
Description=Nemotron Marketing Queue Worker
After=network.target

[Service]
Type=simple
User=renew
WorkingDirectory=/home/renew/workspace/nemotron-mini
ExecStart=/home/renew/workspace/nemotron-mini/venv/bin/python -m app.worker_main --poll-interval-sec 2 --max-jobs-per-tick 1
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
